# gateway/as400_gateway.conf
[program:as400_gateway]
command=/usr/bin/python3 /opt/as400_gateway.py
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/as400_gateway.log
directory=/opt
user=root
startsecs=10
stopasgroup=true
killasgroup=true
